# LPE-CT/installation/debian-trixie/lpe-ct.service
[Unit]
Description=LPE-CT DMZ Management API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lpe-ct
Group=lpe-ct
WorkingDirectory=/opt/lpe-ct/src/LPE-CT
EnvironmentFile=/etc/lpe-ct/lpe-ct.env
ExecStart=/opt/lpe-ct/bin/lpe-ct
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=30
StartLimitBurst=5
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/opt/lpe-ct /var/lib/lpe-ct /var/spool/lpe-ct

[Install]
WantedBy=multi-user.target
