# LPE-CT/installation/debian-trixie/lpe-ct.service
[Unit]
Description=LPE-CT DMZ Management API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__LPE_CT_SERVICE_USER__
Group=__LPE_CT_SERVICE_GROUP__
WorkingDirectory=__LPE_CT_SRC_DIR__/LPE-CT
EnvironmentFile=__LPE_CT_ENV_FILE__
ExecStart=__LPE_CT_BIN_DIR__/lpe-ct
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=30
StartLimitBurst=5
UMask=0027
# Host maintenance actions use sudo for /opt/lpe-ct/bin/lpe-ct-host-action only.
NoNewPrivileges=false
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=false
LockPersonality=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=__LPE_CT_INSTALL_ROOT__ __LPE_CT_STATE_DIR__ __LPE_CT_SPOOL_DIR__ __LPE_CT_LOG_DIR__

[Install]
WantedBy=multi-user.target
